# scripts/etc/udev/rules.d/99-input.rules
KERNEL=="input[0-9]*", RUN+="/usr/local/sbin/a2dp-autoconnect"
